# deployment/ansible-vitam-extra/roles/curator/templates/systemd/vitam-curator-metricbeat-indexes.timer.j2
[Unit]
Description=Delete old topbeat elasticsearch-log indexes

[Timer]
OnCalendar=*-*-* 00:10:00
Persistent=false

[Install]
WantedBy=timers.target
