# sources/metadata/metadata-rest/src/main/systemd/vitam-metadata-purge-dip.timer
[Unit]
Description=Launch vitam dip purge for metadata

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=false


[Install]
WantedBy=timers.target
