# sources/metadata/metadata-rest/src/main/systemd/vitam-metadata-computed-inherited-rules.timer
[Unit]
Description=Launch vitam store graph script for metadata

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=false


[Install]
WantedBy=timers.target
