# sources/ihm-demo/ihm-demo-front/src/deb/control/postinst
#!/bin/sh

systemctl enable vitam-${unix.name}.service >/dev/null 2>&1 || true
