# sources/storage/storage-engine/storage-engine-server/src/deb/control/prerm
#!/bin/sh
case "$1" in
  remove|upgrade|deconfigure)
    # systemd: On uninstall (not upgrade), disable and stop the unit
    if [ -x /bin/systemctl ]; then
      systemctl --no-reload disable vitam-${unix.name}.service >/dev/null 2>&1 || true
      systemctl stop vitam-${unix.name}.service >/dev/null 2>&1 || true

      systemctl --no-reload disable vitam-storageaccess-log-backup.timer >/dev/null 2>&1 || true
      systemctl stop vitam-storage-accesslog-backup.timer >/dev/null 2>&1 || true
      systemctl --no-reload disable vitam-storage-accesslog-backup.service >/dev/null 2>&1 || true
      systemctl stop vitam-storage-accesslog-backup.service >/dev/null 2>&1 || true

      systemctl --no-reload disable vitam-storage-log-backup.timer >/dev/null 2>&1 || true
      systemctl stop vitam-storage-log-backup.timer >/dev/null 2>&1 || true
      systemctl --no-reload disable vitam-storage-log-backup.service >/dev/null 2>&1 || true
      systemctl stop vitam-storage-log-backup.service >/dev/null 2>&1 || true

      systemctl --no-reload disable vitam-storage-log-traceability.timer >/dev/null 2>&1 || true
      systemctl stop vitam-storage-log-traceability.timer >/dev/null 2>&1 || true
      systemctl --no-reload disable vitam-storage-log-traceability.service >/dev/null 2>&1 || true
      systemctl stop vitam-storage-log-traceability.service >/dev/null 2>&1 || true
    fi
    ;;

  failed-upgrade)
    ;;

  *)
    echo "prerm called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac
