# sources/storage/cas-manager/vitam-offer/src/main/systemd/vitam-offer-log-compaction.timer
[Unit]
Description=Vitam offer log compaction timer

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=false

[Install]
WantedBy=timers.target
