# deployment/ansible-vitam/roles/storage/templates/systemd/vitam-storage-backup-logs.timer.j2
[Unit]
Description=Backup old vitam storage logs

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=false

[Install]
WantedBy=timers.target
