# deb/vitam-product/vitam-elasticsearch-head/DEBIAN/postinst
#!/bin/sh

systemctl enable vitam-elasticsearch-data
systemctl start vitam-elasticsearch-data
