# workspace/workspace-rest/src/main/systemd/workspace.service
[Unit]
Description=Workspace Rest
After=network.target
After=syslog.target

[Service]
Type=simple
ExecStart=/vitam/${unix.name}/bin/wrapper.sh

[Install]
WantedBy=multi-user.target
