# sources/storage/storage-engine/storage-engine-server/src/main/systemd/vitam-storage-log-backup.service
[Unit]
Description=Vitam storage log backup

Wants=network.target
After=network.target

[Service]
Type=oneshot
ExecStart=/vitam/script/storage/launchStorageLogBackup.sh
User=vitam
Group=vitam
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
