# sources/storage/storage-engine/storage-engine-server/src/main/systemd/vitam-storage-accesslog-backup.timer
[Unit]
Description=Vitam storage accesslog backup

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=false


[Install]
WantedBy=timers.target
