# sources/storage/storage-engine/storage-engine-server/src/main/systemd/vitam-storage-accesslog-backup.service
[Unit]
Description=Vitam storage accesslog backup

Wants=network.target
After=network.target

[Service]
Type=oneshot
ExecStart=/vitam/script/storage/launchStorageAccessLogBackup.sh
User=vitam
Group=vitam
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target
