# deployment/ansible-vitam/roles/antivirus/templates/clamav/clamav-daemon.socket.j2
[Socket]
ListenStream=/run/clamav/clamd.ctl
ListenStream=127.0.0.1:{{ clamav.port }}
SocketUser=clamav
SocketGroup=clamav
RemoveOnStop=True

[Install]
WantedBy=sockets.target