# deployment/ansible-vitam/roles/curator/templates/systemd/vitam-curator-delete-old-indexes.timer.j2
[Unit]
Description=Delete old elasticsearch-log indexes

[Timer]
OnCalendar=*-*-* 00:20:00
Persistent=false

[Install]
WantedBy=timers.target
