# sources/functional-administration/functional-administration-rest/src/main/systemd/vitam-create-accession-register-symbolic.service
[Unit]
Description=Functional Administrations
After=network.target

[Service]
ExecStart=/vitam/script/functional-administration/vitam-create-accession-register-symbolic.sh
User=vitam
Group=vitam

[Install]
WantedBy=multi-user.target
