# sources/functional-administration/functional-administration-rest/src/main/systemd/vitam-create-accession-register-symbolic.timer
[Unit]
Description=Launch vitam creation of accession register each 6 hours.

[Timer]
OnCalendar=*-*-* 00/6:00:00

[Install]
WantedBy=timers.target
